# conf/crontab
@weekly /home/moztw/conf/backup.sh
@daily  /home/moztw/conf/mirrorlist.sh
